# build-deb/orbit_1.1.3_amd64/lib/systemd/system/orbit.service
[Unit]
Description=Orbit Server Management Panel
After=network.target
Documentation=https://github.com/grosman-net/orbit

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/orbit --config /etc/orbit/config.json
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/orbit

[Install]
WantedBy=multi-user.target

